# duplicati/rootfs/etc/s6-overlay/s6-rc.d/duplicati/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================
# shellcheck shell=bash

bashio::log.info "starting duplicati service"


export _CLI_ARGS=$(bashio::config 'CLI_ARGS')

exec \
    s6-notifyoncheck -d -n 300 -w 1000 -c "nc -z localhost 8200" \
        cd /app/duplicati s6-setuidgid abc mono Duplicati.Server.exe \
        --webservice-interface=any --server-datafolder=/config --webservice-allowed-hostnames=* ${_CLI_ARGS}
