# duplicacy/rootfs/etc/s6-overlay/s6-rc.d/service-duplicacy/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

bashio::log.info "starting Duplicacy server"
umask "${UMASK}"

HOME="${APP_DIR}" exec s6-setuidgid duplicacy "${APP_DIR}/duplicacy_web"
