# duplicacy2/rootfs/etc/s6-overlay/s6-rc.d/service-duplicacy/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

bashio::log.info "starting Duplicacy server"

export APP_FILE_PATH=$(bashio::cache.get 'APP_FILE_PATH')

echo "duplicacy_web running as user $(id -un):$(id -gn)\($(id -u):$(id -g)\) and app file path is is ${APP_FILE_PATH}"

if [ ! -d ~/.duplicacy-web ]; then
    ln -s /config/  ~/.duplicacy-web
fi

touch /logs/duplicacy_web.log

if [ ! -f /config/settings.json ]; then
    echo '{
        "listening_address"     : "0.0.0.0:3875",
        "log_directory"         : "/logs",
        "temporary_directory"   : "/cache"
    }'          > /config/settings.json
fi

if [ ! -f /config/duplicacy.json ]; then
    echo '{}'   > /config/duplicacy.json
fi

chown -R "$USR_ID":"$GRP_ID"   /config /logs /cache

echo "Logging tail of the log from this moment on"
tail -0 -f /logs/duplicacy_web.log &

echo "Starting duplicacy_web"
exec s6-setuidgid ${USR_ID}:${GRP_ID} "${APP_FILE_PATH}"
