# kopia/rootfs/etc/s6-overlay/s6-rc.d/init-kopia-config/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set +e

bashio::log.info "starting kopia config"

if [[ -z "$(bashio::config 'KOPIA_UI_USERNAME')" ]]; then
    echo "ERROR: You must set a password for the Web UI"
    exit 1
fi

SUBJECT="/O=Nginx Proxy/OU=Dummy Certificate/CN=*"
if [[ -f /config/keys/kopia.key && -f /config/keys/kopia.crt ]]; then
    echo "using keys found in /config/keys"
else
    echo "generating self-signed keys in /config/keys, you can replace these with your own keys if required"
    mkdir -p /config/keys
    rm -f \
        /config/keys/kopia.key \
        /config/keys/kopia.crt || true
    openssl req -new -newkey rsa:2048 -x509 -days 3650 -nodes -out /config/keys/kopia.crt -keyout /config/keys/kopia.key -subj "$SUBJECT"
fi
