# borgbackup/rootfs/etc/s6-overlay/s6-rc.d/borgbackup/up
#!/usr/bin/with-contenv bashio

exec /usr/bin/backup.sh
