# traefik/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================
# shellcheck shell=bash

bashio::net.wait_for 443
bashio::log.info "Starting Nginx for dashboard..."

exec nginx -g 'daemon off; pid /tmp/nginx.pid;'
