# traefik/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================
# shellcheck shell=bash

# Wait for traefik to become available
bashio::net.wait_for 8080

bashio::log.info "Starting Nginx ..."
exec nginx
