# traefik/rootfs/etc/s6-overlay/s6-rc.d/init-traefik-config/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

export TZ=$(bashio::config 'timezone')
export log_level=$(bashio::config 'log_level')
export dynamic_configuration_path=$(bashio::config 'dynamic_configuration_path')
export env_vars=$(bashio::config 'env_vars')
export api_user=$(bashio::config 'api_user')
export api_passwd=$(bashio::config 'env_passwd')
echo "
----------------------------------------------------------------------
ENVIRONMENT APP
----------------------------------------------------------------------
TZ=${TZ}
log_level=${log_level}
dynamic_configuration_path=${dynamic_configuration_path}
env_vars=${env_vars}
api_user=${api_user}
api_passwd=${api_passswd}
api_hash_passwd=${api_passswd}
----------------------------------------------------------------------
"
bashio::log.info "Hashing Traefik API password..."
htpasswd -bc /config/traefik/passwd $(bashio::config 'api_user') $(bashio::config 'env_passwd')

bashio::log.info "Ensuring SSL directory..."
mkdir -p /ssl/traefik
bashio::log.info "Ensuring config directory..."
mkdir -p /config/traefik

bashio::log.info "Generating static config..."
gomplate -f /etc/traefik/traefik.yaml.gotmpl -d options=/data/options.json -o /etc/traefik/traefik.yaml
bashio::log.info "Static config generated"
