# traefik/rootfs/etc/s6-overlay/s6-rc.d/init-traefik-config/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

export TZ=$(bashio::config 'timezone')
echo "
----------------------------------------------------------------------
ENVIRONMENT APP
----------------------------------------------------------------------
TZ=${TZ}
log_level=${log_level}
dynamic_configuration_path=${dynamic_configuration_path}
env_vars=${env_vars}
----------------------------------------------------------------------
"
bashio::log.info "Ensuring SSL directory..."
mkdir -p /ssl/traefik
bashio::log.info "Ensuring config directory..."
mkdir -p /config/traefik

bashio::log.info "Generating static config..."
gomplate -f /etc/traefik/traefik.yaml.gotmpl -d options=/data/options.json -o /etc/traefik/traefik.yaml
bashio::log.info "Static config generated"
