# kopia/rootfs/etc/s6-overlay/s6-rc.d/svc-kopia/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

export _timezone=$(bashio::config 'timezone')
export _REFRESH_INTERVAL=$(bashio::config 'REFRESH_INTERVAL')
export _KOPIA_UI_USERNAME=$(bashio::config 'KOPIA_UI_USERNAME')
export _KOPIA_UI_PASSWORD=$(bashio::config 'KOPIA_UI_PASSWORD')
export _keyfile=$(bashio::config 'keyfile')
export _certfile=$(bashio::config 'certfile')

bashio::log.info "starting kopia server"

if bashio::config.true 'ssl'; then
    exec \
        s6-notifyoncheck -d -n 300 -w 1000 \
         s6-setuidgid abc /app/kopia/kopia --timezone ${_timezone} server start \
            --disable-csrf-token-checks \
            --tls-key-file=${_keyfile} \
            --tls-cert-file=${_certfile} \
            --enable-actions \
            --address=0.0.0.0:51515 \
            --refresh-interval=${_REFRESH_INTERVAL} \
            --server-username=${_KOPIA_UI_USERNAME} \
            --server-password=${_KOPIA_UI_PASSWORD}
else
    exec \
        s6-notifyoncheck -d -n 300 -w 1000 \
         s6-setuidgid abc /app/kopia/kopia --timezone ${_timezone} server start \
            --disable-csrf-token-checks \
            --insecure \
            --enable-actions \
            --address=0.0.0.0:51515 \
            --refresh-interval=${_REFRESH_INTERVAL} \
            --server-username=${_KOPIA_UI_USERNAME} \
            --server-password=${_KOPIA_UI_PASSWORD}
fi
