# kopia/rootfs/etc/s6-overlay/s6-rc.d/svc-kopia/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

bashio::log.info "starting kopia server"

if bashio::config.true 'ssl'; then
    exec \
        s6-notifyoncheck -d -n 300 -w 1000 \
        s6-setuidgid abc /config/kopia/kopia --timezone "$(bashio::config 'timzeone') server start \
            --disable-csrf-token-checks \
            --tls-key-file="$(bashio::config 'keyfile')" \
            --tls-cert-file="$(bashio::config 'certfile')" \
            --enable-actions \
            --address=0.0.0.0:51515 \
            --refresh-interval="$(bashio::config 'REFRESH_INTERVAL')" \
            --server-username="$(bashio::config 'KOPIA_UI_USERNAME')" \
            --server-password="$(bashio::config 'KOPIA_UI_PASSWORD')"
else
    exec \
        s6-notifyoncheck -d -n 300 -w 1000 \
        s6-setuidgid abc /config/kopia/kopia --timezone "$(bashio::config 'timzeone') server start \
            --disable-csrf-token-checks \
            --insecure \
            --enable-actions \
            --address=0.0.0.0:51515 \
            --refresh-interval="$(bashio::config 'REFRESH_INTERVAL')" \
            --server-username="$(bashio::config 'KOPIA_UI_USERNAME')" \
            --server-password="$(bashio::config 'KOPIA_UI_PASSWORD')"
fi
