# kopia/rootfs/etc/s6-overlay/s6-rc.d/svc-kopia/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "starting kopia server"

if [[ ${KOPIA_HTTPS} == "true" ]]; then
    exec \
        s6-notifyoncheck -d -n 300 -w 1000 \
        s6-setuidgid abc /app/kopia/kopia server start \
            --disable-csrf-token-checks \
            --tls-key-file=/config/keys/kopia.key \
            --tls-cert-file=/config/keys/kopia.crt \
            --enable-actions \
            --address=0.0.0.0:51515 \
            --refresh-interval="${REFRESH_INTERVAL}" \
            --server-username="${KOPIA_UI_USERNAME}" \
            --server-password="${KOPIA_UI_PASSWORD}"
else
    exec \
        s6-notifyoncheck -d -n 300 -w 1000 \
        s6-setuidgid abc /app/kopia/kopia server start \
            --disable-csrf-token-checks \
            --insecure \
            --enable-actions \
            --address=0.0.0.0:51515 \
            --refresh-interval="${REFRESH_INTERVAL}" \
            --server-username="${KOPIA_UI_USERNAME}" \
            --server-password="${KOPIA_UI_PASSWORD}"
fi
