# traefik/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: traefik
# Configures NGINX for use with the traefik reverse proxy
# ==============================================================================
bashio::log.info "NGINX init config"
declare traefik_port=8080
declare traefik_protocol=http

bashio::log.info "Generate upstream configuration"
bashio::var.json \
    port "^${traefik_port}" \
    | tempio \
        -template /etc/nginx/templates/upstream.gtpl \
        -out /etc/nginx/includes/upstream.conf

bashio::log.info "Generate Ingress configuration"
bashio::var.json \
    interface "$(bashio::addon.ip_address)" \
    port "^$(bashio::addon.ingress_port)" \
    protocol "${traefik_protocol}" \
    | tempio \
        -template /etc/nginx/templates/ingress.gtpl \
        -out /etc/nginx/servers/ingress.conf
